# sdm-ssl-renewal.service
[Unit]
Description=SDM SSL Certificate Renewal Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/bin/ssl_client.py --config /etc/sdm-ssl/config.json
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/ssl/sdmcert /var/log
PrivateTmp=true
